# .config/systemd/user/reflector.service
[Unit]
Description=Weekly Reflector Update

[Service]
Type=oneshot
ExecStart=/usr/bin/reflector --sort rate --save /etc/pacman.d/mirrorlist

